# service.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true

[program:hhvm]
command=/usr/bin/hhvm  -m server -vServer.Type=fastcgi -vServer.Port=9000
autostart=true
autorestart=true
redirect_stderr=true

